# contrib/supervisor/mayan-edms.conf
[program:mayan-edms]
command = /usr/share/mayan-edms/contrib/scripts/gunicorn_start.sh
user = www-data
autostart = true
autorestart = true
redirect_stderr = true

