# contrib/upstart/upstart.conf
start on started rc
stop on stopped rc

respawn
respawn limit 3 5

setuid mayan
setgid mayan

exec /home/mayan/production/deploy/production/service_demon.sh